# install/roles/system/files/rc.local
#!/bin/bash

echo ds2482 0x18 > /sys/bus/i2c/devices/i2c-10/new_device
